# src/Script/job.cron
00 18 * * * python /home/prod/Script/autojob.py
